# deploy/supervisord.ini
[program:bird_personal]
user=root
command=/home/marc/envs/p343-dj1.7/bin/python manage.py runserver 0.0.0.0:8111
directory= /home/marc/python_test_apps/bird/
autostart=true
autorestart=true

stdout_logfile = /tmp/bird_personal.log

redirect_stderr=true
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=50
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
loglevel=info

